# debug/bpf/perf-monitor.service
[Unit]
Description=Performance Monitoring Service
After=network.target

[Service]
ExecStart=/bin/bash /usr/local/bin/perf.sh start
ExecStop=/bin/bash /usr/local/bin/perf.sh stop
Restart=always

[Install]
WantedBy=multi-user.target
